# firmware/systemd/bmw-wifi-apply.path
[Unit]
Description=Watch for BMW Logo Wi-Fi provision request

[Path]
PathExists=/var/lib/bmw-logo/wifi-request.json
PathChanged=/var/lib/bmw-logo/wifi-request.json
Unit=bmw-wifi-apply.service

[Install]
WantedBy=multi-user.target
