# firmware/systemd/bmw-logo-display.service
[Unit]
Description=BMW Logo HDMI Display Renderer (desktop / X11)
After=graphical.target display-manager.service bmw-logo-api.service
Wants=graphical.target bmw-logo-api.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/bmw-logo
Environment=PYTHONPATH=/opt/bmw-logo
Environment=DISPLAY=:0
Environment=SDL_VIDEODRIVER=x11
# Xauthority is rewritten by fix-systemd-paths / setup for the real user
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 2
ExecStart=/opt/bmw-logo/.venv/bin/python -m firmware.display.hdmi_renderer
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
