# firmware/systemd/dot-wifi-watch.service
[Unit]
Description=Dot Wi-Fi watch (auto-join iPhone hotspot)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/dot-wifi-watch
Restart=always
RestartSec=5
Nice=10

[Install]
WantedBy=multi-user.target
