# firmware/systemd/dot-wifi-apply.path
[Unit]
Description=Watch for Dot Wi-Fi provision request

[Path]
# Only on change/create of the request file (avoid double-start with PathExists+PathChanged).
PathChanged=/var/lib/dot/wifi-request.json
Unit=dot-wifi-apply.service

[Install]
WantedBy=multi-user.target
