# firmware/systemd/dot-wifi-boot.service
[Unit]
Description=Dot Wi-Fi first-boot / hotspot join
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/dot-wifi-boot
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
