# firmware/systemd/dot-display-kiosk.service
[Unit]
Description=Dot HDMI Display Renderer (kiosk, no desktop)
After=local-fs.target systemd-udevd.service
Wants=dot-api.service
Before=getty@tty1.service
Conflicts=getty@tty1.service
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
# Root + KMSDRM (no seat). Do not bind a TTY — tty hangup was restarting us in a loop.
User=root
WorkingDirectory=/opt/dot
Environment=PYTHONPATH=/opt/dot
Environment=SDL_VIDEODRIVER=KMSDRM
Environment=SDL_AUDIODRIVER=dummy
Environment=SDL_VIDEO_EGL_DRIVER=libEGL.so.1
Environment=SDL_VIDEO_GL_DRIVER=libGLESv2.so.2
Environment=XDG_RUNTIME_DIR=/run/dot-display
ExecStartPre=/bin/sh -c 'mkdir -p /run/dot-display; chmod 700 /run/dot-display; for i in 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15 16 17 18 19 20; do [ -e /dev/dri/card0 ] || [ -e /dev/dri/card1 ] || [ -e /dev/fb0 ] && exit 0; sleep 0.5; done; exit 0'
ExecStartPre=/bin/sh -c 'if [ -e /dev/fb0 ]; then dd if=/dev/zero of=/dev/fb0 bs=1024 count=4096 status=none 2>/dev/null || true; fi'
StandardInput=null
StandardOutput=journal
StandardError=journal
ExecStart=/opt/dot/.venv/bin/python -m firmware.display.hdmi_renderer
Restart=on-failure
RestartSec=3
TimeoutStopSec=3
KillMode=control-group
FinalKillSignal=SIGKILL
SendSIGKILL=yes

[Install]
WantedBy=multi-user.target
